# recipes-bsp/tdt-driver/tdt-driver/aotom_udev.rules
ACTION=="add", KERNEL=="aotom", SUBSYSTEM=="module", RUN+="/bin/mknod /dev/vfd c 147 0"
